# packages/burlap/burlap-0.9.124.tar.gz/burlap-0.9.124/burlap/templates/openvpn/supervisor.template.conf
[program:{daemon_name}]
environment=
command={command}
directory={directory}
user={user}
process_name=%(program_name)s
stdout_logfile={log_path}
stderr_logfile={log_path}
autostart=true
autorestart=true
startretries=3
priority=998
